# sat-dabctl.conf
[program:sat-0FIP-dabctl]
command=mplayer -quiet -loop 0 -af resample=48000:nowaveheader,format=s16le,channels=2 -ao pcm:file=/home/pi/dab-ctl/conf/audio0.fifo:fast http://odr:odr@localhost:9981/play/stream/service/8a4689ee0ffa82519a68fb7b0052c425
autostart=false
autorestart=true   
stderr_logfile=/var/log/supervisor/sox-fip.err.log                                                                                                                                                     
stdout_logfile=/var/log/supervisor/sox-fip.out.log

[program:sat-1FG-dabctl]
command=mplayer -quiet -loop 0 -af resample=48000:nowaveheader,format=s16le,channels=2 -ao pcm:file=/home/pi/dab-ctl/conf/audio1.fifo:fast http://odr:odr@localhost:9981/play/stream/service/0df64d5d0a90fb8a1ed3dc196280d1c8
autostart=false
autostart=false
autorestart=true   
stderr_logfile=/var/log/supervisor/sox-fg.err.log                                                                                                                                                     
stdout_logfile=/var/log/supervisor/sox-fg.out.log

[program:sat-2OUIFM-dabctl]
command=mplayer -quiet -loop 0 -af resample=48000:nowaveheader,format=s16le,channels=2 -ao pcm:file=/home/pi/dab-ctl/conf/audio2.fifo:fast http://odr:odr@localhost:9981/play/stream/service/4d20746121dc24b2a40b489014d9b74d
autorestart=true   
stderr_logfile=/var/log/supervisor/sox-ouifm.err.log                                                                                                                                                     
stdout_logfile=/var/log/supervisor/sox-ouifm.out.log

[program:sat-3MELODY-dabctl]
command=mplayer -quiet -loop 0 -af resample=48000:nowaveheader,format=s16le,channels=2 -ao pcm:file=/home/pi/dab-ctl/conf/audio3.fifo:fast http://odr:odr@localhost:9981/play/stream/service/f0cb674e8ed586b3217738d23eabfcd0
autostart=false
autorestart=true   
stderr_logfile=/var/log/supervisor/sox-melody.err.log                                                                                                                                                     
stdout_logfile=/var/log/supervisor/sox-melody.out.log

[program:sat-4NOVA-dabctl]
command=mplayer -quiet -loop 0 -af resample=48000:nowaveheader,format=s16le,channels=2 -ao pcm:file=/home/pi/dab-ctl/conf/audio4.fifo:fast http://odr:odr@localhost:9981/play/stream/service/1ec8da0c53d5032dc7f231dcf50326d3
autostart=false
autorestart=true   
stderr_logfile=/var/log/supervisor/sox-nova.err.log                                                                                                                                                     
stdout_logfile=/var/log/supervisor/sox-nova.out.log

[program:sat-5FUN-dabctl]
command=mplayer -quiet -loop 0 -af resample=48000:nowaveheader,format=s16le,channels=2 -ao pcm:file=/home/pi/dab-ctl/conf/audio5.fifo:fast http://odr:odr@localhost:9981/play/stream/service/6a6c7690a8ba3c00fb9fa1a2697e2a94
autostart=false
autorestart=true   
stderr_logfile=/var/log/supervisor/sox-fun.err.log                                                                                                                                                     
stdout_logfile=/var/log/supervisor/sox-fun.out.log

[program:sat-6NOSTALGIE-dabctl]
command=mplayer -quiet -loop 0 -af resample=48000:nowaveheader,format=s16le,channels=2 -ao pcm:file=/home/pi/dab-ctl/conf/audio6.fifo:fast http://odr:odr@localhost:9981/play/stream/service/5c9c907089dc9f815a8cc5535894ddaa
autostart=false
autorestart=true   
stderr_logfile=/var/log/supervisor/sox-fun.err.log                                                                                                                                                     
stdout_logfile=/var/log/supervisor/sox-fun.out.log

[program:sat-7NRJ-dabctl]
command=mplayer -quiet -loop 0 -af resample=48000:nowaveheader,format=s16le,channels=2 -ao pcm:file=/home/pi/dab-ctl/conf/audio7.fifo:fast http://odr:odr@localhost:9981/play/stream/service/b2bb3d744c64094eecdcdbcc5a98754e
autostart=false
autorestart=true   
stderr_logfile=/var/log/supervisor/sox-nrj.err.log                                                                                                                                                     
stdout_logfile=/var/log/supervisor/sox-nrj.out.log
