# config/supervisor/mux.conf
[program:*Multiplex]
command=odr-dabmux /home/pi/dab/config/conf.mux
autostart=false
autorestart=true
stderr_logfile=/var/log/supervisor/mux.log
stdout_logfile=/var/log/supervisor/mux.log
